# systemd/.config/systemd/user/pb-notify.service
[Unit]
Description=Pushbullet notifications

[Service]
ExecStart=/usr/bin/pb-notify
EnvironmentFile=/home/tsani/.pb-notify-env
