# systemd/.config/systemd/user/udiskie.service
[Unit]
Description=Automatic udisks2 mounting

[Service]
ExecStart=/usr/bin/udiskie -a
Type=simple

[Install]
WantedBy=default.target
